# .config/systemd/user/douban.service
[Unit]
Description=A Duoban Crawler

[Service]
Type=oneshot
ExecStart=/home/hill/bin/douban_mpd.sh
;Restart=always
