# .config/systemd/user/indieshuffle.service
[Unit]
Description=A indieshuffle Crawler

[Service]
Type=oneshot
ExecStart=/home/hill/bin/indieshufflstream.py
;Restart=always
