# .config/systemd/user/wallpaper.service
[Unit]
Description=Update wallpaper from unsplash

[Service]
Type=oneshot
ExecStart=/home/hill/bin/unsplashwallpaper
;Restart=always
